# crontab
#run python script every minutes
*/10 * * * *  . /etc/profile.d/rails_env.sh; /bin/sh /app/run-python.sh > /proc/1/fd/1 2>/proc/1/fd/2

